# cookbooks/tor/templates/default/tor-dns-proxy.erb

#!/bin/sh

### BEGIN INIT INFO
# Provides:          tor-dns-proxy
# Required-Start:    $local_fs $remote_fs $network $syslog
# Required-Stop:     $local_fs $remote_fs $network $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: starts the tor-dns-proxy
# Description:       starts tor-dns-proxy using start-stop-daemon
#
# start-stop-daemon --start --make-pidfile --pidfile /var/run/tor-dns-proxy.pid --chuid nobody:nogroup --exec /usr/local/bin/proxymachine -- -p <%= node["tor"]["socks-resolver-port"] %> -c /opt/tor-dns-proxy.rb
### END INIT INFO

PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
DAEMON=/usr/local/bin/proxymachine 
DAEMON_OPTS=-p <%= node["tor"]["socks-resolver-port"] %> -c /opt/tor-dns-proxy.rb
NAME=tor-dns-proxy
DESC=tor-dns-proxy
PID=/var/run/tor-dns-proxy.pid

test -x $DAEMON || exit 0

set -e

. /lib/lsb/init-functions

case "$1" in
  start)
    echo -n "Starting $DESC: "
    start-stop-daemon --start --quiet --make-pidfile --pidfile $PID \
        --exec $DAEMON -- $DAEMON_OPTS || true
    echo "$NAME."
    ;;

  stop)
    echo -n "Stopping $DESC: "
    start-stop-daemon --stop --quiet --pidfile $PID \
        --exec $DAEMON || true
    echo "$NAME."
    ;;

  restart|force-reload)
    echo -n "Restarting $DESC: "
    start-stop-daemon --stop --quiet --pidfile \
        $PID --exec $DAEMON || true
    sleep 1
    start-stop-daemon --start --quiet --pidfile \
        $PID --exec $DAEMON -- $DAEMON_OPTS || true
    echo "$NAME."
    ;;

  reload)
    echo -n "Reloading $DESC configuration: "
    test_nginx_config
    start-stop-daemon --stop --signal HUP --quiet --pidfile $PID \
        --exec $DAEMON || true
    echo "$NAME."
    ;;

  status)
    status_of_proc -p $PID "$DAEMON" nginx && exit 0 || exit $?
    ;;
  *)
    echo "Usage: $NAME {start|stop|restart|reload|force-reload|status}" >&2
    exit 1
    ;;
esac

exit 0
